# deploy/systemd/keba-home-api-reader.service
[Unit]
Description=KEBA Home API Reader
After=network-online.target keba-home-service@carport.service keba-home-service@eingang.service
Wants=network-online.target

[Service]
Type=simple
User=keba
Group=keba
WorkingDirectory=/opt/keba_home_api
EnvironmentFile=/etc/keba/keba-home-api-reader.env
ExecStart=/opt/keba_home_api/keba_api
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/keba

[Install]
WantedBy=multi-user.target
